# files/system/etc/profile.d/settheming.sh
gsettings set org.gnome.desktop.interface gtk-theme 'Adwaita-dark' 
gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'
gsettings set org.gnome.desktop.interface icon-theme 'Papirus'
